# goreplay.ini
[program:gor]
command=/gor --input-raw :%(ENV_APP_PORT) --output-resurface http://localhost:3000/messages
