# docker-resurface/goreplay.ini
[program:gor]
command=/gor --input-raw :%(ENV_APP_PORT)s --input-raw-track-response --output-resurface http://localhost:4001/message --output-resurface-rules "%(ENV_RULES)s"
